# tuersshd/tuersshd.service
[Unit]
Description=Tuer SSH daemon
After=networking.target

[Service]
ExecStart=/usr/bin/tuersshd -privkey=/home/tuersshd/.ssh/id_rsa
User=tuersshd
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
